# profile.d/autojump.sh
../../Cellar/autojump/22.5.3_2/etc/profile.d/autojump.sh